# systemd/.local/share/systemd/user/wallpaper.timer
[Unit]
Description=Update wallpaper every 30min

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
